# dock00-sec/haproxy_watchdog/haproxy-watchdog.service
[Unit]
Description=Watchdog haproxy — redémarre le conteneur si arrêté
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/haproxy-watchdog.sh
