# OPS/_cron/backup-portainer.cron
15 3 * * * /usr/local/bin/backup-portainer.sh >> /var/log/backup-portainer.log 2>&1