# OPS/_cron/fs_poe_off.cron
0 23 * * * /usr/local/sbin/fs_poe_off.sh