# OPS/_cron/lwmd-weekly.cron
0 0 * * 0 /usr/local/sbin/lwmd-weekly.sh >> /var/log/lwmd-weekly.log 2>&1