# OPS/_cron/nina-rclone-sync-s3-to-backup.cron
30 0 * * * /usr/local/sbin/rclone-sync-s3-to-backup.sh