# OPS/_cron/restart-dl-stack-if-unhealthy.cron
*/15 * * * * /usr/local/sbin/restart-dl-stack-if-unhealthy.sh >> /var/log/restart-dl-stack-if-unhealthy.log 2>&1
